# resources/trivy-db-update.service
[Unit]
Description=Update trivy database

[Service]
Type=oneshot
ExecStart=trivy server --download-db-only
